# proj/x86_64/inittab
::sysinit:/etc/init.d/rcS
::respawn:/sbin/getty -L ttyS0 115200 vt100

# Added virtual terminals
# Using mingetty on tty1 and tty2
::respawn:/sbin/mingetty tty1
::respawn:/sbin/mingetty tty2

# Using agetty on tty3 to tty6
::respawn:/sbin/getty -L tty2 38400 vt100
::respawn:/sbin/getty -L tty3 38400 vt100
::respawn:/sbin/getty -L tty4 38400 vt100
::respawn:/sbin/getty -L tty5 38400 vt100
::respawn:/sbin/getty -L tty6 38400 vt100


::shutdown:/bin/umount -a -r
::restart:/sbin/init
